# Oryx.VoteApp.Server/supervisor.conf/voteapp.ini
[program:voteApp]
command=/usr/bin/dotnet /app/voteApp/Oryx.VoteApp.Server.dll
directory=/app/voteApp/
autostart=true
autorestart=true
stderr_logfile=/var/log/VoteApp.err.log
stdout_logfile=/var/log/VoteApp.out.log
environment=HOME=/var/www/,ASPNETCORE_ENVIRONMENT=Production
user=root
stopsignal=INT
stopasgroup=true
killasgroup=true